# launchd/love.citizen-fleet.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>love.citizen-fleet</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/macair/love-unlimited/tools/citizen-fleet.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>LOVE_DIR</key><string>/Users/macair/love-unlimited</string>
    <key>FLEET_PER_TICK</key><string>1</string>
    <key>FLEET_DAILY_CAP</key><string>48</string>
    <key>FLEET_INTERVAL_SEC</key><string>1800</string>
    <key>FLEET_BEAT_BOUND_SEC</key><string>600</string>
    <key>FLEET_AGENTIC_PER_DAY</key><string>12</string>
    <key>FLEET_DAILY_BUDGET_USD</key><string>5.00</string>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/macair/.local/bin</string>
  </dict>
  <key>StartInterval</key><integer>1800</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/macair/love-unlimited/memory/fleet-launchd.log</string>
  <key>StandardErrorPath</key><string>/Users/macair/love-unlimited/memory/fleet-launchd.log</string>
  <key>ThrottleInterval</key><integer>1800</integer>
  <key>ProcessType</key><string>Background</string>
  <key>LowPriorityIO</key><true/>
  <key>Nice</key><integer>5</integer>
</dict>
</plist>
