# tools/love.kosmem.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>love.kosmem</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yu/Love/tools/kosmem-consolidate.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>StandardOutPath</key>
    <string>/Users/yu/Love/memory/.kos/consolidation-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/yu/Love/memory/.kos/consolidation-launchd.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>Nice</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/yu</string>
    </dict>
</dict>
</plist>
