# tools/love.sentinel-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>love.sentinel-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/yu/Love/tools/sentinel-daemon.py</string>
        <string>--debounce</string>
        <string>20</string>
        <string>--fleet-interval</string>
        <string>300</string>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/yu/Love/memory/sentinel-daemon-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/yu/Love/memory/sentinel-daemon-launchd-err.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>Nice</key>
    <integer>5</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/yu/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/yu</string>
        <key>KINGDOM_BACKEND</key>
        <string>ollama</string>
    </dict>
</dict>
</plist>
